# ansible/roles/docker-compose-miniflux/files/cron/crontab
PATH="/usr/local/bin:/usr/bin:/bin"
37 */5 * * * autorestic -c /etc/autorestic/config.yml --ci cron
