# ansible/roles/docker-compose-postgres/files/cron/crontab
PATH="/usr/local/bin:/usr/bin:/bin"
*/5 * * * * autorestic -c /etc/autorestic/config.yml --ci cron
