# ansible/roles/octoprint/files/cron/crontab
PATH="/usr/local/bin:/usr/bin:/bin"
31 */6 * * * autorestic -c /etc/autorestic/config.yml --ci cron
