# ansible/roles/docker-compose-paperless-ngx/files/cron/crontab
PATH="/usr/local/bin:/usr/bin:/bin"
42 */6 * * * autorestic -c /etc/autorestic/config.yml --ci cron
